# curator4-662749731.conf
[program:curator4-662749731]
command=python3 /home/curatorN/curator4-662749731/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curatorN/curator4-662749731/test.err.log
stdout_logfile=/home/curatorN/curator4-662749731/test.out.log